# sub_server.service
[Unit]
Description=Subtitle server
After=network.target

[Service]
Type=simple
ExecStart=/srv/sub_server/sub_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
